# deploy/politrain-pool.timer
[Unit]
Description=Nightly Politrain pool replenishment at 04:30 (after the 03:30 DB backup)

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
